# sol-pricer/deploy/systemd/kalshi-bots.target
[Unit]
Description=All kalshi-pricer / eth-pricer / sol-pricer bots and dashboards
Wants=kalshi-pricer-selective.service
Wants=kalshi-pricer-aggressive.service
Wants=kalshi-pricer-web.service
Wants=eth-pricer-selective.service
Wants=eth-pricer-aggressive.service
Wants=eth-pricer-web.service
Wants=sol-pricer-selective.service
Wants=sol-pricer-aggressive.service
Wants=sol-pricer-web.service

[Install]
WantedBy=default.target
