# sol-pricer/deploy/systemd/sol-pricer-selective.service
[Unit]
Description=sol-pricer SOL selective bot (live trader, no Telegram listener — BTC selective owns it)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/chris/kalshi-pricers/sol-pricer
ExecStart=/home/chris/kalshi-pricers/sol-pricer/.venv/bin/python trade.py --live --profile selective --yes-i-know --no-telegram-listen
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/chris/kalshi-pricers/sol-pricer/logs/selective.log
StandardError=append:/home/chris/kalshi-pricers/sol-pricer/logs/selective.log
TimeoutStopSec=30

[Install]
WantedBy=default.target
