# sol-pricer/deploy/systemd/sol-pricer-aggressive.service
[Unit]
Description=sol-pricer SOL aggressive bot (live trader, no Telegram listener)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/chris/kalshi-pricers/sol-pricer
ExecStart=/home/chris/kalshi-pricers/sol-pricer/.venv/bin/python trade.py --live --profile aggressive --yes-i-know --no-telegram-listen
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/chris/kalshi-pricers/sol-pricer/logs/aggressive.log
StandardError=append:/home/chris/kalshi-pricers/sol-pricer/logs/aggressive.log
TimeoutStopSec=30

[Install]
WantedBy=default.target
