# deploy/systemd/kalshi-pricer-selective.service
[Unit]
Description=kalshi-pricer BTC selective bot (live trader, owns Telegram listener)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/chris/kalshi-pricer
ExecStart=/home/chris/kalshi-pricer/.venv/bin/python trade.py --live --profile selective --yes-i-know
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/chris/kalshi-pricer/logs/selective.log
StandardError=append:/home/chris/kalshi-pricer/logs/selective.log
TimeoutStopSec=30

[Install]
WantedBy=default.target
